# scripts/inkctl.service
[Unit]
Description=k9s InkCTL Monitor
# Wait for xochitl to be running before we start monitoring it
After=xochitl.service

[Service]
# Lightweight monitor mode: no Qt, just watches for the trigger notebook.
# When triggered, starts inkctl-display.service which conflicts
# with xochitl so systemd handles the display handoff cleanly.
ExecStart=/home/root/inkctl --monitor
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
